# systemd/dead-mans-switch.timer
[Unit]
Description=Dead man switch timer

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
